# supervisord/supervisord.conf
[unix_http_server]
file=/var/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[supervisorctl]
serverurl=unix:///var/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:test]
user=test
directory=/opt/test/bin
command=/opt/test/bin/test
autostart=true
autorestart=false
startsecs=1
startretries=999
stderr_logfile=/var/log/test/test_err.log
stdout_logfile=/var/log/test/test_out.log
